# stridesync/rootfs/etc/services.d/mcp-server/run
#!/usr/bin/with-contenv bashio
# s6 service: serves synced activity data over MCP (Streamable HTTP via mcp-proxy).
# See app/mcp/server.py and PROJECT_PLAN.md §1 (MCP server) / §2 (MCP Connection).

bashio::log.info "Starting StrideSync MCP server..."

cd /app
exec python3 -m app.mcp.server
